# supervisor_gunicorn.ini
[program:suntv] 
command=/opt/suntv/.venv/bin/gunicorn -k gevent -w 2 -b 127.0.0.1:5000 manage:app 
directory=/opt/suntv
user=root      
autostart=true 
autorestart=true 
stdout_logfile=/var/log/supervisor/suntv.log
stderr_logfile=/var/log/supervisor/suntv-err.log
startsecs=0                                                                             ; 启动时间
stopwaitsecs=0                                                                          ; 终止等待时间